# develop/docker/files/fms/systemd/first-start.service
[Unit]
Description=First start service, only executes once on first container start.

[Service]
Type=simple
ExecStart=-/first-start
