# develop/docker/files/fms/systemd/install_certificate.service
[Unit]
Description=Install certificate into FileMaker Server

[Service]
Type=simple
ExecStart=-/install_certificate.sh
